# hosts/instance1/k3s-agents/media-worker/media-worker-ensure.service
[Unit]
Description=Ensure the K3s media-worker starts after delayed storage
Documentation=https://docs.k3s.io/
Requires=k3s.service
After=k3s.service network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=-/usr/bin/systemctl reset-failed media-worker.service home-bupd-hdd-data.mount
ExecStart=/usr/bin/systemctl start media-worker.service
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
