# experiments/alpine-bootc-pi5/system_files/etc/init.d/cloudflared
#!/sbin/openrc-run

name=cloudflared
command="/usr/bin/cloudflared"
command_args="${cloudflared_args:-tunnel run}"
command_background=true
pidfile="/run/cloudflared.pid"
output_log="/var/log/cloudflared.log"
error_log="/var/log/cloudflared.err"

depend() {
    need networking
}
