# hosts/homelab/k3s-agents/media-worker/media-worker.container
[Unit]
Description=K3s media-worker agent
Documentation=https://docs.k3s.io/
Wants=network-online.target
After=network-online.target k3s.service
Requires=k3s.service

[Container]
Image=docker.io/rancher/k3s:v1.36.2-k3s1
Pull=never
ContainerName=media-worker
HostName=media-worker

Network=media-worker.network
IP=10.89.0.2

# K3s runs its own containerd and CNI inside this system container. NVIDIA CDI
# gives the worker container the physical GPU; the in-cluster GPU operator is
# responsible for exposing it to application Pods.
PodmanArgs=--privileged --device=nvidia.com/gpu=all --security-opt=label=disable

Volume=/etc/rancher/k3s-media-worker/config.yaml:/etc/rancher/k3s/config.yaml:ro
Volume=/etc/rancher/k3s-media-worker/token:/run/secrets/k3s-agent-token:ro
Volume=/etc/rancher/k3s-media-worker/node:/etc/rancher/node
Volume=/var/lib/rancher/k3s-media-worker:/var/lib/rancher/k3s
Volume=/var/lib/rancher/k3s-media-worker-kubelet:/var/lib/kubelet:rshared
Volume=/var/lib/rancher/k3s-media-worker-cni:/var/lib/cni
Volume=/var/log/k3s-media-worker:/var/log
Volume=/home/bupd/hdd/data:/home/bupd/hdd/data:rw
Volume=/lib/modules:/lib/modules:ro
Volume=/etc/cdi:/etc/cdi:ro
# The host's CDI spec injects dynamically linked NVIDIA tools and libraries.
# Supply the glibc runtime that the intentionally tiny K3s image does not carry,
# so GPU Operator can validate the pre-installed host driver.
Volume=/usr/lib/ld-linux-x86-64.so.2:/lib64/ld-linux-x86-64.so.2:ro
Volume=/usr/lib/libc.so.6:/usr/lib/libc.so.6:ro
Volume=/usr/lib/libpthread.so.0:/usr/lib/libpthread.so.0:ro
Volume=/usr/lib/libm.so.6:/usr/lib/libm.so.6:ro
Volume=/usr/lib/libdl.so.2:/usr/lib/libdl.so.2:ro
Volume=/usr/lib/librt.so.1:/usr/lib/librt.so.1:ro

Exec=agent

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
